# ansible/roles/jetson_nano_services/files/rover-jetson-odrive-back.service
[Unit]
<<<<<<< HEAD:ansible/roles/jetson_nano_services/files/rover-jetson-odrive-back.service
Description=MRover Drive
=======
Description=MRover Jetson Science Bridge
>>>>>>> main:ansible/roles/jetson_services/files/rover-jetson-science-bridge.service

[Service]
EnvironmentFile=/etc/default/rover.env
User={{ running_user }}
<<<<<<< HEAD:ansible/roles/jetson_nano_services/files/rover-jetson-odrive-back.service
ExecStart={{ workspace_dir }}/jarvis exec jetson_odrive_bridge 2
=======
ExecStart={{ workspace_dir }}/jarvis exec jetson_science_bridge
>>>>>>> main:ansible/roles/jetson_services/files/rover-jetson-science-bridge.service
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
