# ansible/roles/onboard_services/files/rover-onboard-teleop.service
[Unit]
Description=MRover Onboard Teleop

[Service]
EnvironmentFile=/etc/default/rover.env
User={{ running_user }}
ExecStart={{ workspace_dir }}/jarvis exec onboard_teleop
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
