# ansible/roles/science_services/files/rover-lcm-bridge-server.service
[Unit]
Description=MRover LCM Bridge Server
>>>>>>> travis-free:ansible/roles/science_services/files/rover-lcm-bridge-server.service

[Service]
User={{ running_user }}
EnvironmentFile=/etc/default/rover.env
ExecStart={{ workspace_dir }}/jarvis exec lcm_bridge_server
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
