# ansible/roles/base_station_services/files/rover-base_station_cameras.service
[Unit]
Description=Base Station Cameras Stream

[Service]
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/mrover/.Xauthority"
User={{ user }}
ExecStart={{ workspace_dir }}/base_station/cameras/run.sh
Restart=always
