# ansible/roles/base_station_services/files/rover-base_station_gstreamer_2.service
[Unit]
Description=Base Station Video Stream 2

[Service]
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/mrover/.Xauthority"
User={{ user }}
ExecStart=/usr/bin/gst-launch-1.0 udpsrc port=5001 caps="application/x-rtp,encoding-name=H264,payload=96" ! rtph264depay ! h264parse ! decodebin ! autovideosink
Restart=always
