# ansible/roles/base_station_services/files/rover-lcm_bridge_server.service
[Unit]
Description=MRover LCM Bridge Server

[Service]
User={{ user }}
EnvironmentFile=/etc/default/rover.env
ExecStart={{ workspace_dir }}/jarvis exec lcm_bridge_server

[Install]
WantedBy=multi-user.target
